# teamd-lvl1.target
[Unit]
Description=Create team interfaces (lvl 1)
Before=network.service

[Install]
WantedBy=multi-user.target
